# config/rc.local
#!/bin/bash

export PATH="/opt/software/bin:$PATH"

ntpdate -u > /dev/null 2>&1 < /dev/null &
xlight -reset &
(rmdir /media/*/* /media/* && rm /media/*) &

if [ "`diff /boot/grub/grub.cfg /boot/grub/grub.cfg-ok`" ]
then
    cp -p /boot/grub/grub.cfg-ok /boot/grub/grub.cfg
fi

if [ -f /etc/rc.mount ]
then
    . /etc/rc.mount
fi

exit 0
